# crontab
# start sync procedure every sunday at 2:00am
0 2 * * SUN ~/sync/scripts/invoke_sync.sh > /dev/null 2>&1
